# makeup/server.conf
start on startup
exec forever start /root/krepstroy/makeup/server.js
